# keytap.service
[Unit]
Description=keytap

[Service]
ExecStart=/usr/local/sbin/keytap local

[Install]
WantedBy=default.target
