# keytap.service
[Unit]
Description=keytap

[Service]
ExecStart=/usr/local/sbin/keytap

[Install]
WantedBy=default.target
